# crontabs/root
*/1 * * * * date "+%Y-%m-%d %H:%M:%S" >> /var/log/cron.log
